# domssr.conf
[program:domssr]
directory=/var/www/domhack.gtfb.tech
command=node server-dist/index.js
numprocs=1
autostart=true
autorestart=true
user=www-data